# roles/timelapse_store/templates/systemd/timelapse_movie.timer
[Unit]
Description=Timelapse Movie Creation timer

[Timer]
OnCalendar=*-*-* 00:58:00

[Install]
WantedBy=timers.target
