# roles/timelapse/templates/systemd/timelapse.service
[Unit]
Description=Timelapse capture service

[Service]
Type=oneshot
ExecStart=/bin/bash {{timelapse_executable}}
