# roles/timelapse/templates/systemd/timelapse.timer
[Unit]
Description=Timelapse image capture

[Timer]
OnUnitActiveSec=300s
OnBootSec=300s

[Install]
WantedBy=timers.target
