# roles/flaky-network/files/systemd/network-restart.timer
[Unit]
Description=Check and restart networking on connection loss

[Timer]
OnUnitActiveSec=60s
OnBootSec=60s

[Install]
WantedBy=timers.target
