# roles/flaky-network/files/systemd/network-restart.service
[Unit]
Description=Check and restart networking on connection loss

[Service]
Type=oneshot
ExecStart=/bin/bash /usr/local/bin/network-restart.sh
